# deploy/crmlog-scheduler.cron
# CRMLog scheduler — install with:
#   sudo cp deploy/crmlog-scheduler.cron /etc/cron.d/crmlog-scheduler
#   sudo chmod 644 /etc/cron.d/crmlog-scheduler
#
# Or append to root crontab:
#   * * * * * cd /var/www/crmlog && php artisan schedule:run >> /dev/null 2>&1

* * * * * www-data cd /var/www/crmlog && php artisan schedule:run >> /dev/null 2>&1
