# deploy/crmlog-queue.service
[Unit]
Description=CRMLog Queue Worker
After=network.target mysql.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
Restart=always
RestartSec=5
WorkingDirectory=/var/www/crmlog
ExecStart=/usr/bin/php /var/www/crmlog/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
Nice=10

[Install]
WantedBy=multi-user.target
